# linux/idlekiller/idle-killer.conf
description "The Grim Reaper for idle user."

start on runlevel [2345]
stop on runlevel [016]

respawn

exec /usr/bin/idle-killer.sh


